# apj-auto/setup/jp.co.elcommun.apj-order.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>jp.co.elcommun.apj-order</string>

  <!-- __APJ_DIR__ は setup_mac.sh がインストール時に実パスへ置換する -->
  <key>ProgramArguments</key>
  <array>
    <string>__APJ_DIR__/venv/bin/python3</string>
    <string>__APJ_DIR__/run_apj_order.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__APJ_DIR__</string>

  <!-- 平日(月〜金) 15:00 に起動。
       スリープ中に発火時刻を過ぎた場合は復帰時に1回実行される(launchd仕様)。
       同日2回目はスクリプト側の state ガードでスキップされる。 -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>__HOME__/apj-order/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/apj-order/logs/launchd.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>LANG</key>
    <string>ja_JP.UTF-8</string>
  </dict>
</dict>
</plist>
